# system_configs/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=false

[program:uwsgi]
command=/usr/bin/uwsgi --die-on-term --ini /etc/uwsgi/vassals/cherrypy.ini
